# templates/etc/systemd/system/kube-scheduler.service.j2
[Unit]
Description=Kubernetes Scheduler
Documentation=https://github.com/GoogleCloudPlatform/kubernetes

[Service]
ExecStart={{k8s_bin_dir}}/kube-scheduler \
  --leader-elect=true \
  --master=http://{{hostvars[inventory_hostname]['ansible_' + etcd_interface].ipv4.address}}:8080 \
  --v=2
Restart=on-failure
RestartSec=5

[Install]
WantedBy=multi-user.target'
